# ignis-server-installer/ignis-server-izpack-installer/src/main/installer/bin/fcr-engine-ignis-server.service
[Unit]
Description=FCR Engine ignis-server ${APP_VER} service
After=sshd.service fcr-engine-platform-tools.service

[Service]
Type=simple
User=${USER_NAME}
WorkingDirectory=${INSTALL_PATH}
PIDFile=${tmp.path}/ignis-server.pid

ExecStart=${INSTALL_PATH}/bin/ignis.sh start
ExecStop=${INSTALL_PATH}/bin/ignis.sh stop
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target