# ignis-design/ignis-design-server/src/assembly/bin/design-studio.service.template
[Unit]
Description=FCR Engine Design Studio service

[Service]
Type=simple
User=${USER_NAME}
WorkingDirectory=${DESIGN_STUDIO_HOME}
PIDFile=/tmp/design-studio-${USER_NAME}.pid

ExecStart=${DESIGN_STUDIO_BIN}/design-studio.sh start
ExecStop=${DESIGN_STUDIO_BIN}/design-studio.sh stop
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target